# config/crontab
# Test: Enable internet access every minute
0 5 * * * /app/disable-internet.sh 192.168.178.10 on

# Monitor cron daemon and restart if necessary
* * * * * pgrep -f crond > /dev/null || { echo "Cron daemon not running, restarting..."; pkill crond; crond; }
